# systemd/intelliops-beat.service
[Unit]
Description=IntelliOps Celery Beat Scheduler
After=network.target redis.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/intelliops
ExecStart=/usr/local/bin/celery -A api.celery_app beat --loglevel=info
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
